# deploy/systemd/spotify-soundeo-sync.service
[Unit]
Description=Spotify to Soundeo sync service
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/opt/spotify-soundeo-sync
EnvironmentFile=/opt/spotify-soundeo-sync/.env
ExecStart=/opt/spotify-soundeo-sync/.venv/bin/python -m app daily-sync
Restart=no

