# deploy/systemd/spotify-soundeo-sync.service
[Unit]
Description=Spotify to Soundeo sync service
After=network-online.target
Wants=network-online.target
Requires=docker.service
After=docker.service

[Service]
Type=oneshot
WorkingDirectory=/opt/spotify-soundeo-sync
ExecStart=/usr/bin/docker compose run --rm app python -m app daily-sync
Restart=no
