# deploy/systemd/spotify-soundeo-sync.timer
[Unit]
Description=Twice-weekly Spotify to Soundeo sync timer

[Timer]
OnCalendar=Tue,Fri *-*-* 19:00:00
Persistent=true

[Install]
WantedBy=timers.target
