# deploy/cron/spotify-soundeo-sync.cron
0 19 * * 2,5 cd /opt/spotify-soundeo-sync && docker compose run --rm app python -m app daily-sync >> logs/cron.log 2>&1
